# docker/runit/service/nginx/run
#!/bin/bash

export LD_LIBRARY_PATH=""

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf

