# docker/runit/service/spark/run
#!/bin/bash
set -e
set -x

exec 2>&1

export APPLICATION_WEB_PROXY_BASE="${DISPATCHER_UI_WEB_PROXY_BASE}"

cd /opt/spark/dist

if [ "$FRAMEWORK_NAME" != "spark" ]; then
    export SPARK_DAEMON_JAVA_OPTS="$SPARK_DAEMON_JAVA_OPTS -Dspark.deploy.zookeeper.dir=/spark_mesos_dispatcher_$FRAMEWORK_NAME"
fi

if [ "$SPARK_DISPATCHER_MESOS_ROLE" != "" ]; then
    export SPARK_DAEMON_JAVA_OPTS="$SPARK_DAEMON_JAVA_OPTS -Dspark.mesos.role=$SPARK_DISPATCHER_MESOS_ROLE"
fi

HISTORY_SERVER_CONF=""
if [ "${ENABLE_HISTORY_SERVER:=false}" = "true" ]; then
    HISTORY_SERVER_CONF="spark.mesos.historyServer.url=${HISTORY_SERVER_WEB_PROXY_BASE}"
fi

sed "s,<HISTORY_SERVER_CONF>,${HISTORY_SERVER_CONF}," \
    conf/mesos-cluster-dispatcher.properties.template >conf/mesos-cluster-dispatcher.properties

function add_if_non_empty() {
	if [ -n "$2" ]; then
		echo "$1=$2" >> conf/mesos-cluster-dispatcher.properties
	fi
}

add_if_non_empty spark.ssl.enabled "${SPARK_SSL_ENABLED}"
add_if_non_empty spark.ssl.keyStoreBase64 "${SPARK_SSL_KEYSTOREBASE64}"
add_if_non_empty spark.ssl.trustStoreBase64 "${SPARK_SSL_TRUSTSTOREBASE64}"
add_if_non_empty spark.ssl.keyPassword "${SPARK_SSL_KEYPASSWORD}"
add_if_non_empty spark.ssl.keyStorePassword "${SPARK_SSL_KEYSTOREPASSWORD}"
add_if_non_empty spark.ssl.trustStorePassword "${SPARK_SSL_TRUSTSTOREPASSWORD}"
add_if_non_empty spark.ssl.protocol "${SPARK_SSL_PROTOCOL}"
add_if_non_empty spark.ssl.enabledAlgorithms "${SPARK_SSL_ENABLEDALGORITHMS}"

exec /opt/spark/dist/bin/spark-class \
    org.apache.spark.deploy.mesos.MesosClusterDispatcher \
    --port "${DISPATCHER_PORT}" \
    --webui-port "${DISPATCHER_UI_PORT}" \
    --master "${MESOS_MASTER}" \
    --zk "${ZK}" \
    --host "${HOST}" \
    --name "${FRAMEWORK_NAME}" \
    --properties-file "conf/mesos-cluster-dispatcher.properties"
